# resources/badboys.service
[Unit]
Description=BAckup DataBases On Your Server

[Service]
Type=oneshot
ExecStart=/usr/local/bin/badboys
